# openrem/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn openremproject.wsgi:application --bind 0.0.0.0:8000 --timeout 300

[program:consumer]
command=python manage.py run_huey
